# deploy/tournaments-bot.service
[Unit]
Description=Padel Tournaments Telegram Bot
After=network.target

[Service]
Type=simple
User=bot
WorkingDirectory=/opt/tournament_bot
ExecStart=/opt/tournament_bot/venv/bin/python bot.py
Restart=always
RestartSec=10
EnvironmentFile=/opt/tournament_bot/.env

[Install]
WantedBy=multi-user.target
